# deploy/systemd/mainlogic-worker-location.service
[Unit]
Description=Main Logic location-updated worker
After=network-online.target rabbitmq-server.service
Requires=rabbitmq-server.service

[Service]
Type=simple
User=mainlogic
Group=mainlogic
WorkingDirectory=/opt/mainlogic
ExecStart=/usr/bin/node dist/workers/locationUpdatedWorker.js
Restart=always
RestartSec=2
Environment=NODE_ENV=production
EnvironmentFile=-/etc/mainlogic/env

[Install]
WantedBy=multi-user.target
